# fileadmin/cron/crontab
0-55/5 * * * *		/usr/bin/php -f $HOME/htdocs/cms/typo3/cli_dispatch.phpsh scheduler