# fileadmin/cron/crontab
*/5 * * * *		/usr/bin/php -f /srv/www/vlib1/htdocs/typo3/cli_dispatch.phpsh scheduler